# tools/install/services/vrg-network.target
[Unit]
Description=Videoreg Is Available To Start Network Services
After=basic.target
PartOf=vrg-boot.target