# tools/install/services/vrg-charging.target
[Unit]
Description=Videoreg Is Charging Now
After=basic.target
PartOf=vrg-boot.target